# ops/systemd/rpow-healthcheck.timer
[Unit]
Description=Probe rpow-server /ready every 90 seconds

[Timer]
OnBootSec=60
OnUnitActiveSec=90
AccuracySec=5

[Install]
WantedBy=timers.target
